# scripts/cron.d/ktp-spike-digest-daily
# KTP daily spike-digest — Tier 3 Project 3 phase 3
#
# Fires daily at 05:00 ET (after the 04:30 ET ktp-perf-rollup, so today's
# digest reads yesterday's complete fleet data — both rollups operate on
# the same trailing-day window). Reads real per-day counts from
# hlstatsx.ktp_spike_daily (plus ktp_spike_signatures for new/gone-quiet
# fingerprints) and posts a Discord embed summarizing:
#
#   - Spikes ≥10ms: yesterday's count per fingerprint vs its trailing-7-day
#     average, with the worst-hit instance
#   - New fingerprints first_seen yesterday (cross-checks the alert hook)
#   - Noise floor: one-line sub-10ms total
#   - Gone quiet: historically-active fingerprints that stopped firing
#
# Default channel: 1497957091107668070 (#ktp-crashes — same as
# PERF_ALERT_CHANNEL + SPIKE_ALERT_CHANNEL). Override via
# SPIKE_DIGEST_CHANNEL in /etc/ktp/discord-relay.conf.
#
# Implementation:
#   /usr/local/bin/ktp-spike-digest    — Python script (this repo: scripts/ktp-spike-digest.py)
#   hlstatsx.ktp_spike_daily           — populated by KTPProfileAggregator daemon (1.5.31)
#   hlstatsx.ktp_spike_signatures      — populated by KTPProfileAggregator daemon
#
# Operator note: ktp_telemetry needs SELECT on ktp_spike_daily in addition
# to the 1.5.20-era ktp_spike_signatures grant — grants are PER-TABLE on
# this DB (see KTPProfileAggregator README.md § "Required MySQL grants").

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=''

# 05:00 ET daily. Reuses the ktp-profile-aggregator venv (pymysql lives
# there; sharing avoids maintaining a second venv just for this script).
0 5 * * *   root  set -a; source /opt/ktp-profile-aggregator/.env; set +a; /opt/ktp-profile-aggregator/venv/bin/python /usr/local/bin/ktp-spike-digest --config /etc/ktp/discord-relay.conf >> /var/log/ktp-spike-digest.log 2>&1
