# data-server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

# ============================================
# MySQL — HLStatsX database
# ============================================
[program:mysql]
command=/usr/sbin/mysqld --user=mysql
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/mysql.log
stderr_logfile=/var/log/supervisor/mysql-err.log

# ============================================
# HLTV Proxy 1 — connects to game server 1
# ============================================
# Config mounted at /opt/hltv/instance-1/hltv.cfg
[program:hltv-1]
command=/opt/hltv/instance-1/hltv
directory=/opt/hltv/instance-1
environment=LD_LIBRARY_PATH="/opt/hltv"
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/supervisor/hltv-1.log
stderr_logfile=/var/log/supervisor/hltv-1-err.log

# ============================================
# HLTV Proxy 2 — connects to game server 2
# ============================================
# Config mounted at /opt/hltv/instance-2/hltv.cfg
[program:hltv-2]
command=/opt/hltv/instance-2/hltv
directory=/opt/hltv/instance-2
environment=LD_LIBRARY_PATH="/opt/hltv"
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/supervisor/hltv-2.log
stderr_logfile=/var/log/supervisor/hltv-2-err.log

# ============================================
# HLStatsX log receiver (stub)
# ============================================
# Listens on UDP 27500 for game server log lines.
# Replace with real HLStatsX Perl daemon when ready.
[program:hlstatsx-stub]
command=socat -u UDP4-LISTEN:27500,fork,reuseaddr STDOUT
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/hlstatsx.log
stderr_logfile=/var/log/supervisor/hlstatsx-err.log

# ============================================
# Node.js backend — HUD Observer
# ============================================
[program:backend]
command=/usr/local/bin/node /app/node_modules/.bin/ts-node --script-mode /app/backend/src/app.ts
directory=/app
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend-err.log

# ============================================
# React frontend — static build served on :3000
# ============================================
[program:frontend]
command=serve -s /app/web/build -l 3000
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend-err.log
